# deploy/clawnode.service
[Unit]
Description=ClawOps Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/clawnode run --config /etc/clawnode/config.json
Restart=always
RestartSec=5
Environment=RUST_LOG=info
WorkingDirectory=/var/lib/clawnode

[Install]
WantedBy=multi-user.target
